# _release/cacophonator-rs485-controller.service
[Unit]
Description=Cacophonator management interface
After=network.target

[Service]
ExecStart=/usr/bin/rs485-controller
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
