# _release/trap-sequence.service
[Unit]
Description=Cacophonator trap sequence
After=network.target

[Service]
ExecStart=/usr/bin/trap-sequence
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
